# Daily_crawler/daily_crawler.cron
0 */3 * * * /home/admin/janesi-spider/Daily_crawler/start_crawl.sh;
*/20 * * * * /home/admin/janesi-spider/Daily_crawler/start_hot_crawl.sh

